# tests/integration_tests/reference_files/celery_sso.supervisor.celery-beat.conf
[program:celery-beat]
user=cloudron
environment=HOME=/home/cloudron,USER=cloudron
directory=/app/code
command=/app/code/venv/bin/celery -A blog beat --loglevel=info --schedule /app/data/celerybeat-schedule
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
